# roles/setdefaultpath/template/etc/profile.d/path.sh
export PATH=/opt/jdk1.8.0_121/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin

export JAVA_HOME=/opt/jdk1.8.0_121
export CLASSPATH="/opt/jdk1.8.0_121/lib:."
