# roles/profile-env-alias-scripts/template/etc/profile.d/git-tips.sh
#!/bin/bash

function lazygit() {
  git add .
  git commit -a -m $*
  git push
}

alias gp=lazygit


