# roles/nodeangular/template/etc/profile.d/pathnode.sh
export PATH=$PATH:/opt/node-v6.11.5-linux-x64/bin/

